# deploy/systemd/todoalarm.service
[Unit]
Description=TaskFlow Server
Documentation=https://github.com/youruser/todoalarm
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=todoalarm
Group=todoalarm
WorkingDirectory=/opt/todoalarm
ExecStart=/opt/todoalarm/todoalarm-server -config /opt/todoalarm/config.toml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=todoalarm

# === 安全加固 ===
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# 只允许写这两个目录
ReadWritePaths=/opt/todoalarm/data /opt/todoalarm/backup

# 资源限制(SQLite + 调度器轻量;给充分余量)
LimitNOFILE=4096
MemoryMax=512M

[Install]
WantedBy=multi-user.target
